# deploy/clara-hub.service
[Unit]
Description=ClaraAFJV - Hub Emploi 3D & Jeu Video (Next.js)
After=network.target postgresql.service

[Service]
Type=simple
User=clara
Group=clara
WorkingDirectory=/home/clara/ClaraAFJV
Environment=NODE_ENV=production
Environment=PATH=/usr/bin:/usr/local/bin:/bin
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clara-hub

[Install]
WantedBy=multi-user.target
